# conf/node_exporter.service
[Unit]
Description=Prometheus Node Exporter (host metrics)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=node-exporter
Group=node-exporter
# systemd collector: nginx, nginx-exporter, node_exporter 서비스 생존 여부를 node_systemd_unit_state 로 노출
ExecStart=/opt/node_exporter/node_exporter --web.listen-address=:9100 --collector.systemd --collector.systemd.unit-include="(nginx|nginx-exporter|node_exporter)\\.service"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
